# .config/systemd/system/mount-gdrive/mount-gdrive.service
[Unit]
Description=Mount Google Drive to predictable paths
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mount-gdrive
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
